# provision/files/80-vagrant-mounted.rules
SUBSYSTEM=="bdi",ACTION=="add",RUN+="/bin/bash /root/.udev-vagrant-mounted.sh"
